# doemaarwatt/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

exec nginx -c /etc/nginx/nginx.conf -e /var/log/nginx/error.log